# supervisord.conf
[supervisord]
nodaemon=true
[program:neauHelper]
command=python /usr/src/app/index.py
autorstart=true
